# extra/redbot.service
[Unit]
Description=REDbot
After=network.target

[Service]
Type=simple
# Static system user shared with redbot-gc.service, so the cleanup timer can
# manage files the daemon writes under StateDirectory. DynamicUser is avoided
# deliberately: it hands each unit a different ephemeral UID, which can't share
# ownership of persistent state. Create the user with the shipped sysusers
# snippet (extra/redbot-sysusers.conf): `systemd-sysusers`.
User=redbot
Group=redbot

# Process
Environment=PYTHONUNBUFFERED='true' SYSTEMD_WATCHDOG='true'
WorkingDirectory=/opt/redbot
ExecStart=/opt/redbot/app/bin/redbot_daemon /opt/redbot/config/config.txt
SyslogIdentifier=redbot
Restart=always
SystemCallArchitectures=native

# Watchdog. The daemon pings every 3s from the event loop, so this doubles as a
# liveness check on the loop itself. Keep it well above the ping interval: at 10s
# an ordinary reclaim stall was enough to trip it.
WatchdogSec=30

# Sandbox
NotifyAccess=main
NoNewPrivileges=true
LockPersonality=true
ProtectSystem=strict
ProtectProc=invisible
ProtectHostname=true
ProtectHome=tmpfs
ProtectClock=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProcSubset=pid
RestrictSUIDSGID=true
RestrictRealtime=true
KeyringMode=private
UMask=0077
PrivateDevices=true
# Preserve the isolation DynamicUser used to imply.
PrivateTmp=true
RemoveIPC=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
CapabilityBoundingSet=
DevicePolicy=closed
RestrictNamespaces=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
StateDirectory=redbot
LogsDirectory=redbot
LimitCORE=0

# Resource Limits. Steady-state working set is ~57M (interpreter, Jinja
# templates, compiled regexes, babel's CLDR tables); MemoryHigh must stay clear
# of it, or the kernel throttles the process in reclaim and the resulting
# multi-second stalls trip the watchdog. Swap is disabled deliberately: with it
# enabled, an overshoot became tens of seconds of thrashing instead of a prompt
# OOM kill and restart. Running with --debug adds ~45M of tracemalloc overhead
# and needs MemoryHigh raised to match.
CPUQuota=60%
MemoryLow=40M
MemoryHigh=96M
MemoryMax=128M
MemorySwapMax=0

[Install]
WantedBy=multi-user.target
